# setup/front/cms/crontab
*/10 2 * * * [[CMS]]/Ideal/Library/sitemap/index.php
*/10 3 * * * [[CMS]]/Ideal/Library/FileMonitor/scan.php